# ops/systemd/trader-watchdog-daily.timer
[Unit]
Description=Trader daily health report at 13:00 UTC

[Timer]
OnCalendar=*-*-* 13:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
