# ops/systemd/trader-watchdog-daily.service
[Unit]
Description=Trader daily health report
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/trader
ExecStart=/var/www/trader/.venv/bin/python /var/www/trader/ops/watchdog.py --daily
